# debian/rabbitmq-consume.target
[Unit]
Description=Description=RabbitMQ Consumer for Virtual Hosts
After=network.target

[Install]
WantedBy=multi-user.target
